# gitlab/overlay/etc/network/if-up.d/gitlab
#!/bin/sh -e

#Reconfigure ip
/usr/lib/inithooks/firstboot.d/35set-gitlab-ip
